# .crontab
# Monsuke Fed Rates
00 08 * * * python /home/eclipse/monsuke/monsuke.py
